# 10-connected-tty.rules
SUBSYSTEM=="tty",ACTION=="add",MODE="0666",RUN+="/home/alina.gotovtceva/udev-intro/hello.sh %k"
