# cp2102.rules
ATTRS{idProduct}=="ea60", ATTRS{idVendor}=="10c4", ACTION=="add", RUN+="/home/alina.gotovtceva/udev-intro/cp2102.py"

